# latest/windows_builtin/win_scheduled_task_deletion.yml
detect:
  case sensitive: false
  log type: wel
  op: is
  path: Event/System/EventID
  target: log
  value: '4699'
respond:
- action: report
  metadata:
    author: David Strassegger
    description: Detects scheduled task deletion events. Scheduled tasks are likely
      to be deleted if not used for persistence. Malicious Software often creates
      tasks directly under the root node e.g. \TASKNAME
    falsepositives:
    - Software installation
    level: medium
    references:
    - https://twitter.com/matthewdunwoody/status/1352356685982146562
    - https://docs.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4699
    tags:
    - attack.execution
    - attack.privilege_escalation
    - car.2013-08-001
    - attack.t1053.005
  name: Scheduled Task Deletion

