# latest/windows_process_creation/proc_creation_win_susp_vbscript_unc2452.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: Execute
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: CreateObject
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: RegRead
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: window.close
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Microsoft\Windows\CurrentVersion
    - case sensitive: false
      not: true
      op: contains
      path: event/COMMAND_LINE
      value: \Software\Microsoft\Windows\CurrentVersion\Run
respond:
- action: report
  metadata:
    author: Florian Roth
    description: Detects suspicious inline VBScript keywords as used by UNC2452
    falsepositives:
    - Unknown
    level: high
    references:
    - https://www.microsoft.com/security/blog/2021/03/04/goldmax-goldfinder-sibot-analyzing-nobelium-malware/
    tags:
    - attack.persistence
    - attack.t1547.001
  name: Suspicious VBScript UN2452 Pattern

