# latest/windows_process_creation/proc_creation_win_apt_actinium_persistence.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: schtasks
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: create
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: wscript
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: ' /e:vbscript'
respond:
- action: report
  metadata:
    author: Andreas Hunkeler (@Karneades)
    description: Detects specific process parameters as used by ACTINIUM scheduled
      task persistence creation.
    falsepositives:
    - Unlikely
    level: high
    references:
    - https://www.microsoft.com/security/blog/2022/02/04/actinium-targets-ukrainian-organizations
    tags:
    - attack.persistence
    - attack.t1053
    - attack.t1053.005
  name: Potential ACTINIUM Persistence Activiy

