# latest/windows_process_creation/proc_creation_win_hktl_powersploit_empire_default_schtasks.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: or
      rules:
      - case sensitive: false
        op: ends with
        path: event/PARENT/FILE_PATH
        value: \powershell.exe
      - case sensitive: false
        op: ends with
        path: event/PARENT/FILE_PATH
        value: \pwsh.exe
    - case sensitive: false
      op: ends with
      path: event/FILE_PATH
      value: \schtasks.exe
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /Create
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: powershell.exe -NonI
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /TN Updater /TR
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /SC ONLOGON
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /SC DAILY /ST
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /SC ONIDLE
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /SC HOURLY
respond:
- action: report
  metadata:
    author: Markus Neis, @Karneades
    description: Detects the creation of a schtask via PowerSploit or Empire Default
      Configuration.
    falsepositives:
    - Unlikely
    level: high
    references:
    - https://github.com/0xdeadbeefJERKY/PowerSploit/blob/8690399ef70d2cad10213575ac67e8fa90ddf7c3/Persistence/Persistence.psm1
    - https://github.com/EmpireProject/Empire/blob/08cbd274bef78243d7a8ed6443b8364acd1fc48b/lib/modules/powershell/persistence/userland/schtasks.py
    tags:
    - attack.execution
    - attack.persistence
    - attack.privilege_escalation
    - attack.s0111
    - attack.g0022
    - attack.g0060
    - car.2013-08-001
    - attack.t1053.005
    - attack.t1059.001
  name: HackTool - Default PowerSploit/Empire Scheduled Task Creation

