# latest/windows_process_creation/proc_creation_win_win10_sched_task_0day.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: ends with
      path: event/FILE_PATH
      value: \schtasks.exe
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /change
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /TN
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /RU
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /RP
respond:
- action: report
  metadata:
    author: Olaf Hartong
    description: Detects Task Scheduler .job import arbitrary DACL write\par
    falsepositives:
    - Unknown
    level: high
    references:
    - https://github.com/SandboxEscaper/polarbearrepo/tree/master/bearlpe
    tags:
    - attack.privilege_escalation
    - attack.t1053.005
    - car.2013-08-001
  name: Windows 10 Scheduled Task SandboxEscaper 0-day

