# latest/windows_process_creation/proc_creation_win_susp_new_service_creation.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: or
      rules:
      - op: and
        rules:
        - case sensitive: false
          op: ends with
          path: event/FILE_PATH
          value: \sc.exe
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: create
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: binPath=
      - op: and
        rules:
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: New-Service
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: -BinaryPathName
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: powershell
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: mshta
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: wscript
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: cscript
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: svchost
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: dllhost
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'cmd '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: cmd.exe /c
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: rundll32
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: C:\Users\Public
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Downloads\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Desktop\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Microsoft\Windows\Start Menu\Programs\Startup\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: C:\Windows\TEMP\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \AppData\Local\Temp
respond:
- action: report
  metadata:
    author: Nasreddine Bencherchali
    description: Detects creation of a new service via "sc" command or the powershell
      "new-service" cmdlet with suspicious binary paths
    falsepositives:
    - Unlikely
    level: high
    references:
    - https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.003/T1543.003.md
    - https://web.archive.org/web/20180331144337/https://www.fireeye.com/blog/threat-research/2018/03/sanny-malware-delivery-method-updated-in-recently-observed-attacks.html
    tags:
    - attack.persistence
    - attack.privilege_escalation
    - attack.t1543.003
  name: Suspicious New Service Creation

