# latest/windows_process_creation/proc_creation_win_reg_add_run_key.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: reg
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: ' ADD '
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: Software\Microsoft\Windows\CurrentVersion\Run
respond:
- action: report
  metadata:
    author: Florian Roth
    description: Detects suspicious command line reg.exe tool adding key to RUN key
      in Registry
    falsepositives:
    - Unknown
    level: medium
    references:
    - https://app.any.run/tasks/9c0f37bc-867a-4314-b685-e101566766d7/
    - https://docs.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys
    tags:
    - attack.persistence
    - attack.t1547.001
  name: Reg Add RUN Key

