# latest/windows_process_creation/proc_creation_win_susp_schtasks_env_folder.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: or
      rules:
      - op: and
        rules:
        - case sensitive: false
          op: ends with
          path: event/FILE_PATH
          value: \schtasks.exe
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' /create '
        - op: or
          rules:
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: '%AppData%'
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: \AppData\Local\
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: \AppData\Roaming\
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: '%Public%'
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: \Users\Public
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: C:\Windows\Temp
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: C:\Perflogs
      - op: and
        rules:
        - case sensitive: false
          op: ends with
          path: event/PARENT/COMMAND_LINE
          value: \svchost.exe -k netsvcs -p -s Schedule
        - op: or
          rules:
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: '%Public%'
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: \Users\Public
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: C:\Windows\Temp
          - case sensitive: false
            op: contains
            path: event/COMMAND_LINE
            value: C:\Perflogs
    - not: true
      op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: update_task.xml
      - case sensitive: false
        op: contains
        path: event/PARENT/COMMAND_LINE
        value: unattended.ini
respond:
- action: report
  metadata:
    author: Florian Roth
    description: Detects Schtask creations that point to a folder references in environment
      variables or often used by malware
    falsepositives:
    - Benign scheduled tasks creations or executions that happen often during software
      installations
    - Software that uses the AppData folder and scheduled tasks to update the software
      in the AppData folders
    level: high
    references:
    - https://www.welivesecurity.com/2022/01/18/donot-go-do-not-respawn/
    - https://www.joesandbox.com/analysis/514608/0/html#324415FF7D8324231381BAD48A052F85DF04
    tags:
    - attack.execution
    - attack.t1053.005
  name: Suspicius Schtasks From Env Var Folder

