# latest/windows_process_creation/proc_creation_win_susp_schtasks_change.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: ends with
      path: event/FILE_PATH
      value: \schtasks.exe
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: ' /Change '
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: ' /TN '
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \AppData\Local\Temp
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \AppData\Roaming\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Users\Public\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \WINDOWS\Temp\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Desktop\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Downloads\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \Temporary Internet
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: C:\ProgramData\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: C:\Perflogs\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '%ProgramData%'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '%appdata%'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '%comspec%'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '%localappdata%'
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: regsvr32
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: rundll32
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'cmd /c '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'cmd /k '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'cmd.exe /c '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'cmd.exe /k '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: powershell
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: mshta
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: wscript
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: cscript
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: certutil
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: bitsadmin
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: bash.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'bash '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: scrcons
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'wmic '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: wmic.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: forfiles
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: scriptrunner
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: hh.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'hh '
respond:
- action: report
  metadata:
    author: Nasreddine Bencherchali
    description: 'Detects when an attacker tries to modify an already existing scheduled
      tasks to run from a suspicious location

      Attackers can create a simple looking task in order to avoid detection on creation
      as it''s often the most focused on

      Instead they modify the task after creation to include their malicious payload

      '
    falsepositives:
    - Unknown
    level: high
    references:
    - Internal Research
    - https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks
    tags:
    - attack.execution
    - attack.t1053.005
  name: Suspicious Modification Of Scheduled Tasks

