# latest/windows_process_creation/proc_creation_win_schtasks_appdata_local_system.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - case sensitive: false
        op: ends with
        path: event/FILE_PATH
        value: \schtasks.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /Create
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /RU
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /TR
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: C:\Users\
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \AppData\Local\
      - op: or
        rules:
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: NT AUT
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' SYSTEM '
    - not: true
      op: and
      rules:
      - op: and
        rules:
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: \AppData\Local\Temp\
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: TeamViewer_.exe
      - case sensitive: false
        op: ends with
        path: event/FILE_PATH
        value: \schtasks.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: /TN TVInstallRestore
respond:
- action: report
  metadata:
    author: pH-T, Nasreddine Bencherchali
    description: Detects the creation of a schtask that executes a file from C:\Users\<USER>\AppData\Local
    falsepositives:
    - Unknown
    level: high
    references:
    - https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
    tags:
    - attack.execution
    - attack.persistence
    - attack.t1053.005
    - attack.t1059.001
  name: Suspicious Schtasks Execution AppData Folder

