# latest/windows_process_creation/proc_creation_win_susp_schtasks_user_temp.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - case sensitive: false
        op: ends with
        path: event/FILE_PATH
        value: \schtasks.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '/Create '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \AppData\Local\Temp
    - not: true
      op: and
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '/Create /TN "klcp_update" /XML '
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \klcp_update_task.xml
respond:
- action: report
  metadata:
    author: frack113
    description: schtasks.exe create task from user AppData\Local\Temp
    falsepositives:
    - Unknown
    level: high
    references:
    - malware analyse https://www.joesandbox.com/analysis/514608/0/html#324415FF7D8324231381BAD48A052F85DF04
    tags:
    - attack.execution
    - attack.t1053.005
  name: Suspicious Add Scheduled Task From User AppData Temp

