# latest/windows_process_creation/win_susp_direct_asep_reg_keys_modification.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: ends with
      path: event/FILE_PATH
      value: \reg.exe
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: add
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\Run
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\RunOnce
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\RunOnceEx
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\RunServices
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\RunServicesOnce
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows NT\CurrentVersion\Windows
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \software\Microsoft\Windows\CurrentVersion\Explorer\User Shell Folders
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: \system\CurrentControlSet\Control\SafeBoot\AlternateShell
respond:
- action: report
  metadata:
    author: Victor Sergeev, Daniil Yugoslavskiy, oscd.community
    description: Detects direct modification of autostart extensibility point (ASEP)
      in registry using reg.exe.
    falsepositives:
    - Legitimate software automatically (mostly, during installation) sets up autorun
      keys for legitimate reasons.
    - Legitimate administrator sets up autorun keys for legitimate reasons.
    level: medium
    references:
    - https://github.com/redcanaryco/atomic-red-team/blob/master/atomics/T1547.001/T1547.001.md
    tags:
    - attack.persistence
    - attack.t1547.001
  name: Direct Autorun Keys Modification

