# experimental/windows_process_creation/proc_creation_win_schtasks_schedule_type.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - op: or
        rules:
        - case sensitive: false
          op: ends with
          path: event/FILE_PATH
          value: \schtasks.exe
        - case sensitive: false
          op: is
          path: event/ORIGINAL_FILE_NAME
          value: schtasks.exe
      - op: or
        rules:
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' ONLOGON '
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' ONSTART '
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' ONCE '
        - case sensitive: false
          op: contains
          path: event/COMMAND_LINE
          value: ' ONIDLE '
    - not: true
      op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: NT AUT
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: ' SYSTEM'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: HIGHEST
respond:
- action: report
  metadata:
    author: Nasreddine Bencherchali (Nextron Systems)
    description: Detects scheduled task creations or modification on a suspicious
      schedule type
    falsepositives:
    - Legitimate processes that run at logon. Filter according to your environment
    level: high
    references:
    - https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change
    - https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create
    - http://blog.talosintelligence.com/2022/09/lazarus-three-rats.html
    tags:
    - attack.execution
    - attack.t1053.005
  name: Suspicious Schtasks Schedule Types

