# experimental/windows_process_creation/proc_creation_win_schtasks_creation.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - case sensitive: false
        op: ends with
        path: event/FILE_PATH
        value: \schtasks.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: ' /create '
    - not: true
      op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/USER_NAME
        value: AUTHORI
      - case sensitive: false
        op: contains
        path: event/USER_NAME
        value: AUTORI
respond:
- action: report
  metadata:
    author: Florian Roth (Nextron Systems)
    description: Detects the creation of scheduled tasks by user accounts via the
      "schtasks" utility.
    falsepositives:
    - Administrative activity
    - Software installation
    level: low
    references:
    - https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create
    tags:
    - attack.execution
    - attack.persistence
    - attack.privilege-escalation
    - attack.t1053.005
    - attack.s0111
    - car.2013-08-001
    - stp.1u
  name: Scheduled Task Creation Via Schtasks.EXE

