# experimental/windows_process_creation/proc_creation_win_schtasks_parent.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - op: and
      rules:
      - case sensitive: false
        op: ends with
        path: event/FILE_PATH
        value: \schtasks.exe
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: '/Create '
      - op: or
        rules:
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: \AppData\Local\
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: \AppData\Roaming\
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: \Temporary Internet
        - case sensitive: false
          op: contains
          path: event/PARENT/FILE_PATH
          value: \Users\Public\
    - not: true
      op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: update_task.xml
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: unattended.ini
respond:
- action: report
  metadata:
    author: Florian Roth (Nextron Systems)
    description: Detects suspicious scheduled task creations from a parent stored
      in a temporary folder
    falsepositives:
    - Software installers that run from temporary folders and also install scheduled
      tasks
    level: medium
    references:
    - https://app.any.run/tasks/649e7b46-9bec-4d05-98a5-dfa9a13eaae5/
    tags:
    - attack.execution
    - attack.t1053.005
  name: Suspicious Add Scheduled Task Parent

