# experimental/windows_process_creation/proc_creation_win_schtasks_reg_loader.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: ends with
      path: event/FILE_PATH
      value: \schtasks.exe
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /Create
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: /SC
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: FromBase64String
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: Get-ItemProperty
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'HKCU:'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'HKLM:'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: 'registry::'
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: HKEY_
respond:
- action: report
  metadata:
    author: '@Kostastsale, @TheDFIRReport, slightly modified by pH-T (Nextron Systems)'
    description: Detects the creation of a schtask that executes a base64 encoded
      payload stored in the Windows Registry using PowerShell.
    falsepositives:
    - Unknown
    level: high
    references:
    - https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
    tags:
    - attack.execution
    - attack.persistence
    - attack.t1053.005
    - attack.t1059.001
  name: Scheduled Task Executing Powershell Encoded Payload from Registry

