# experimental/macos_process_creation/proc_creation_macos_persistence_via_plistbuddy.yml
detect:
  events:
  - NEW_PROCESS
  - EXISTING_PROCESS
  op: and
  rules:
  - op: is mac
  - op: and
    rules:
    - case sensitive: true
      op: ends with
      path: event/FILE_PATH
      value: /PlistBuddy
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: RunAtLoad
    - case sensitive: false
      op: contains
      path: event/COMMAND_LINE
      value: 'true'
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: LaunchAgents
      - case sensitive: false
        op: contains
        path: event/COMMAND_LINE
        value: LaunchDaemons
respond:
- action: report
  metadata:
    author: Sohan G (D4rkCiph3r)
    description: Detects potential persistence activity using LaunchAgents or LaunchDaemons
      via the PlistBuddy utility
    falsepositives:
    - Unknown
    level: high
    references:
    - https://redcanary.com/blog/clipping-silver-sparrows-wings/
    - https://www.manpagez.com/man/8/PlistBuddy/
    tags:
    - attack.persistence
    - attack.t1543.001
    - attack.t1543.004
  name: Potential Persistence Via PlistBuddy

