# experimental/windows_builtin/win_security_susp_scheduled_task_update-edr.yml
detect:
  event: WEL
  op: and
  rules:
  - op: is windows
  - op: and
    rules:
    - case sensitive: false
      op: is
      path: event/EVENT/System/EventID
      value: '4702'
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \AppData\Local\Temp\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \AppData\Roaming\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \Users\Public\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \WINDOWS\Temp\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: C:\Temp\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \Desktop\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \Downloads\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: \Temporary Internet
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: C:\ProgramData\
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: C:\Perflogs\
    - op: or
      rules:
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: regsvr32
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: rundll32
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: cmd.exe</Command>
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: cmd</Command>
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: '<Arguments>/c '
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: '<Arguments>/k '
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: '<Arguments>/r '
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: powershell
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: pwsh
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: mshta
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: wscript
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: cscript
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: certutil
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: bitsadmin
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: bash.exe
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: 'bash '
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: scrcons
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: 'wmic '
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: wmic.exe
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: forfiles
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: scriptrunner
      - case sensitive: false
        op: contains
        path: event/EVENT/EventData/TaskContentNew
        value: hh.exe
respond:
- action: report
  metadata:
    author: Nasreddine Bencherchali (Nextron Systems)
    description: Detects update to a scheduled task event that contain suspicious
      keywords.
    falsepositives:
    - Unknown
    level: high
    references:
    - https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4698
    tags:
    - attack.execution
    - attack.privilege_escalation
    - attack.persistence
    - attack.t1053.005
  name: Suspicious Scheduled Task Update

